# dotfiles/dot-config/systemd/system/update-texlive.service
[Unit]
Description=Update texlive using tlmgr

[Service]
Type=oneshot
ExecStart=/usr/local/bin/tlmgr update --all

[Install]
WantedBy=default.target
