# dotfiles/dot-config/systemd/user/update-cargo.timer
[Unit]
Description=Update user cargos daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
